# roles/service-nginx-ffapi/templates/ffapi-sync.timer.j2
#
# {{ ansible_managed }}
#
[Unit]
Description=Sync {{ item.name }} Freifunk API repo timer

[Timer]
OnBootSec=5m
OnUnitActiveSec=1d

[Install]
WantedBy=timers.target
