# roles/service-bird-icvpn/templates/icvpn-tinc-bgp-update.service.j2
#
# {{ ansible_managed }}
#
[Unit]
Description=Update tinc hosts and bgp peers for icvpn

[Service]
ExecStart=/home/admin/clones/backend-scripts/gen_icvpn_bgp_gw.sh
User=admin
Group=admin
