# roles/service-bird-icvpn/templates/icvpn-update.timer.j2
#
# {{ ansible_managed }}
#
[Unit]
Description=Timer which schedules icvpn-update.service

[Timer]
OnBootSec=1h
OnUnitActiveSec=2d

[Install]
WantedBy=timers.target
