# roles/service-fastd-mesh/templates/fastd-peer-limit-update.service.j2
#
# {{ ansible_managed }}
#
[Unit]
Description=Update fastd peer limits

[Service]
ExecStart=/home/admin/clones/backend-scripts/limit_fastd_peers_gw.py
User=admin
Group=admin
