# roles/service-fastd-mesh/templates/fastd-sync-meshkeys.service.j2
#
# {{ ansible_managed }}
#
[Unit]
Description=Update fastd mesh peers repos

[Service]
ExecStart=/home/admin/clones/backend-scripts/sync_meshkeys_gw.sh
User=admin
Group=admin
