# roles/service-fastd-mesh/templates/fastd-peer-limit-update.timer.j2
#
# {{ ansible_managed }}
#
[Unit]
Description=Timer which schedules fastd-peer-limit-update.service

[Timer]
OnBootSec=5min
OnUnitActiveSec=5min

[Install]
WantedBy=timers.target
