# roles/service-fastd-mesh/templates/fastd-limiter-limit.service.j2
#
# {{ ansible_managed }}
#
[Unit]
Description=Update fastd peer limits

[Service]
Type=oneshot
ExecStart={{ fastd_limiter_path }} limit
User=admin
Group=admin
