# roles/service-fastd-mesh/templates/fastd-limiter-limit.timer.j2
#
# {{ ansible_managed }}
#
[Unit]
Description=Timer which schedules fastd-limiter-limit.service

[Timer]
OnBootSec=5min
OnUnitActiveSec=5min

[Install]
WantedBy=timers.target
