# roles/service-fastd-mesh/templates/fastd-limiter-peers.service.j2
#
# {{ ansible_managed }}
#
[Unit]
Description=Update fastd connected peers

[Service]
Type=oneshot
ExecStart={{ fastd_limiter_path }} peers
User=admin
Group=admin
