# roles/service-fastd-mesh/templates/fastd-limiter-peers.timer.j2
#
# {{ ansible_managed }}
#
[Unit]
Description=Timer which schedules fastd-limiter-peers.service

[Timer]
OnBootSec=5min
OnUnitActiveSec=1min

[Install]
WantedBy=timers.target
