# roles/service-fastd-mesh/templates/fastd-sync-meshkeys.timer.j2
#
# {{ ansible_managed }}
#
[Unit]
Description=Update fastd mesh peers repos timer

[Timer]
OnBootSec=5m
OnUnitActiveSec=15m

[Install]
WantedBy=timers.target
