# roles/service-bind-slave/templates/icvpn-dns-update.timer.j2
#
# {{ ansible_managed }}
#
[Unit]
Description=Timer which schedules icvpn-bind-update.service

[Timer]
OnBootSec=2h
OnUnitActiveSec=2d

[Install]
WantedBy=timers.target
