# roles/service-bundle-mail/templates/etc/systemd/system/timer-modoboa-statsupdate.timer.j2
#
# {{ ansible_managed }}
#
[Unit]
Description=Timer which schedules timer-modoboa-statsupdate.service

[Timer]
OnBootSec=45m
OnUnitActiveSec=1h

[Install]
WantedBy=timers.target
