# roles/service-bundle-mail/templates/etc/systemd/system/timer-modoboa-logs-table-cleanup.timer.j2
#
# {{ ansible_managed }}
#
[Unit]
Description=Timer which schedules timer-modoboa-logs-table-cleanup.service

[Timer]
OnCalendar= *-*-* *:00:00 

[Install]
WantedBy=timers.target
